# systemd/backup-cleanup.timer
[Unit]
Description=Daily trigger for the backup / cleanup service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target